# etc/cron.d/ing-sync
# Sync ing files every 5 minutes
*/5 * * * * root /etc/threat-trace/ing-sync-files
